# etc/supervisor/conf.d/chrome.conf
[program:chrome]
command=/opt/google/chrome/google-chrome --no-first-run --force-device-scale-factor=2 --remote-debugging-port=19222 --enable-benchmarking --enable-net-benchmarkin
priority=10
user=automation
directory=/home/automation
environment=HOME="/home/automation"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_backups=4
stdout_logfile_maxbytes=10MB
stderr_logfile=NONE
